# src/etc/systemd/system/pi3ctrl-core.service
[Unit]
Description=Pi3Ctrl core agent
After=multi-user.target
 
[Service]
Type=simple
User=root
Group=www-data
WorkingDirectory=/opt/pi3ctrl/
Environment="VIRTUAL_ENV=/opt/py3"
Environment="PATH=$VIRTUAL_ENV/bin:$PATH"
ExecStart=
ExecStart=$VIRTUAL_ENV/bin/pi3ctrl-core
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi3ctrl-core
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
