# src/etc/systemd/system/pi3ctrl-http.service
[Unit]
Description=Pi3Ctrl http agent
After=network.target
 
[Service]
Type=simple
User=smh
Group=www-data
WorkingDirectory=/opt/pi3ctrl
Environment="VIRTUAL_ENV=/opt/py3"
Environment="PATH=$VIRTUAL_ENV/bin:$PATH"
ExecStart=
ExecStart=$VIRTUAL_ENV/bin/uwsgi --ini $VIRTUAL_ENV/share/pi3ctrl/uwsgi.ini
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi3ctrl-http
Restart=always
RestartSec=30
 
[Install]
WantedBy=multi-user.target
